# addon/autoclear/auto-clear-ram.service
[Unit]
Description=Auto Clearing RAM
Documentation=https://JsPhantom.com
After=syslog.target network-online.target

[Service]
User=root
NoNewPrivileges=true
ExecStart=/usr/bin/auto-clear-ram
Restart=on-failure

[Install]
WantedBy=multi-user.target
